# wifi-sleep.timer
[Unit]
Description=Suspend computer at 22:00 daily

[Timer]
OnCalendar=*-*-* 22:00:00
Persistent=true

[Install]
WantedBy=timers.target
